# services/launchd/bet.schnapp.failover-freshness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>bet.schnapp.failover-freshness</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/schnapp/code/schnapp-bet/services/launchd/op-wrap.sh</string>
        <string>/usr/local/bin/python3</string>
        <string>/Users/schnapp/code/schnapp-bet/services/failover/freshness_check.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/schnapp/code/schnapp-bet/services/failover</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/schnapp/code/schnapp-bet/services/failover/freshness-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/schnapp/code/schnapp-bet/services/failover/freshness-check.err.log</string>
</dict>
</plist>
